# conf/supervise-jenkinspot.conf
[program:jenkinspot]
command=/opt/jenkinspot/jenkinspot.py
redirect_stderr=true
stdout_logfile=/opt/jenkinspot/logs/jenkinspot.log
autorestart=true