# conf/supervise-jenkinspot.conf
[program:jenkinspot]
command=/opt/jenkinspot/jenkinspot.py
redirect_stderr=true
stdout_logfile=/opt/jenkinspot/logs/access.log
autorestart=true